# heartbeat-agent/deploy/com.quorbz.nexus-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Nexus Heartbeat Agent — launchd plist (macOS agents)
     Install to: ~/Library/LaunchAgents/com.quorbz.nexus-heartbeat.plist
     Load: launchctl load ~/Library/LaunchAgents/com.quorbz.nexus-heartbeat.plist
     Replace {{NODE_PATH}} and {{HOME}} before installing. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quorbz.nexus-heartbeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{NODE_PATH}}</string>
        <string>{{HOME}}/nexus-heartbeat/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{HOME}}/nexus-heartbeat</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/nexus-heartbeat/heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/nexus-heartbeat/heartbeat.error.log</string>
</dict>
</plist>
